# .config/autostart/chrome.desktop
[Desktop Entry]
Exec=google-chrome
Name=Chrome
Type=Application
