# .config/autostart/chrome.desktop
[Desktop Entry]
Exec=google-chrome
Name=Firefox
Type=Application
