# .config/autostart/firefox.desktop
[Desktop Entry]
Exec=firefox-esr
Name=Firefox
Type=Application
