# .config/autostart/skypeforlinux.desktop
[Desktop Entry]
Exec=/usr/bin/skypeforlinux
Name=Skype
Type=Application
