# .config/autostart/telegram.desktop
[Desktop Entry]
Exec=telegram-desktop
Name=Telegram
Type=Application
