# .config/autostart/terminator.desktop
[Desktop Entry]
Exec=terminator
Name=Terminator
Type=Application
