# .config/autostart/Yandex.Disk.desktop
[Desktop Entry]
Exec=yandex-disk start
Name=Yandex.Disk
NoDisplay=true
Type=Application
